# systemd/system/hugetlbfs.service.j2
[Unit]
Description=Set up hugetlbfs for readout

[Service]
Type=oneshot
ExecStart={{ o2_install_path }}/bin/setup_hugetlbfs.sh

[Install]
WantedBy=multi-user.target
